# files/changeme.service
{{ ansible_managed | comment }}

[Unit]
Description=foobar service
Wants=foobar

[Service]
Type=single
User=root
ExecStart=/bin/bash
